# site-cookbooks/andare/files/default/andare.conf
description "cpprefjp Updating Server"
start on runlevel [2345];
stop on runlevel [!2345];
respawn
exec su andare -c 'bash /home/andare/andare/runserver.sh'
