# site-cookbooks/wandbox/files/default/kennel2.conf
description "Frontend of Sandbox Compiler"
start on runlevel [2345]
stop on runlevel [!2345]
respawn
exec su wandbox -c '/usr/local/kennel2/bin/kennel --syslog -c /usr/local/kennel2/etc/kennel.json'
