# site-cookbooks/ghwebhook/files/default/ghwebhook.conf
description "Github Webhook Server"
start on runlevel [2345]
stop on runlevel [!2345]
respawn
chdir /home/ghwebhook/ghwebhook
exec su ghwebhook -c 'python ghwebhook.py'
